# foldingathome/rootfs/etc/s6-overlay/s6-rc.d/fah/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# shellcheck disable=SC2191
# ==============================================================================
# Home Assistant Community Add-on: Folding@home
# Runs the FAHClient
# ==============================================================================
declare -a options

bashio::log.info 'Starting Folding@home...'

ACCOUNT_TOKEN="$(bashio::config 'account_token')"
MACHINE_NAME="$(bashio::config 'machine_name')"

FAH_USER="$(bashio::config 'user')"
FAH_TEAM="$(bashio::config 'team')"
FAH_PASSKEY="$(bashio::config 'passkey')"

options+=(--http-addresses 0.0.0.0:7396)
options+=(--allow 0/0)

options+=(--user "${FAH_USER}")
options+=(--team "${FAH_TEAM}")

# Passkey is optional
if [[ -n "${FAH_PASSKEY}" ]]; then
  options+=(--passkey "${FAH_PASSKEY}")
fi


cd /data || bashio::exit.nok "Failed to change to working directory"

CONFIG_FILE="/data/config.xml"

xml_escape() {
  local value="$1"
  value=${value//&/&amp;}
  value=${value//</&lt;}
  value=${value//>/&gt;}
  value=${value//\"/&quot;}
  value=${value//\'/&apos;}
  printf '%s' "${value}"
}

# Only write config.xml when linking to an account (v8-5 web control machine linking).
# For anonymous/default runs, do not force a config.xml (v8.5.5 does not support fold-anon here).
if [[ -n "${ACCOUNT_TOKEN}" ]]; then
  bashio::log.info "Writing ${CONFIG_FILE} (account linked)"
  {
    echo "<config>"
    echo "  <account-token v=\"$(xml_escape "${ACCOUNT_TOKEN}")\"/>"
    echo "  <machine-name v=\"$(xml_escape "${MACHINE_NAME}")\"/>"
    echo "</config>"
  } > "${CONFIG_FILE}"
  chmod 600 "${CONFIG_FILE}"
else
  rm -f "${CONFIG_FILE}"
fi

# Log args, but redact passkey to avoid leaking secrets into add-on logs
log_args=()
skip_next=false
for arg in "${options[@]}"; do
  if ${skip_next}; then
    log_args+=("********")
    skip_next=false
    continue
  fi

  if [[ "${arg}" == "--passkey" ]]; then
    log_args+=("--passkey")
    skip_next=true
    continue
  fi

  log_args+=("${arg}")
done

bashio::log.info "FAHClient args: ${log_args[*]}"
exec /usr/bin/FAHClient "${options[@]}"
